# conf/rapiddiskd.service
[Unit]
Description=A rapiddisk network management daemon
After=network.target

[Service]
ExecStart=/sbin/rapiddiskd
Type=forking
Restart=always

[Install]
WantedBy=default.target
